# home 3/api-server.service
[Unit]
Description=Camera API Server
After=network.target picamera-wifi-bootstrap.service picamera-time-sync.service
Wants=network.target picamera-time-sync.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
RuntimeDirectory=picamera
RuntimeDirectoryMode=0755
MemoryMax=96M
MemoryHigh=80M
CPUWeight=80
Nice=5
# wifi_mode は上書きしない（setdefault のみ）。毎回 ap 固定するとテザリングと競合してフラップしやすい。
ExecStartPre=/usr/bin/python3 -c "import json,os; p='/home/pi/camera_settings.json'; d={}; exec(\"try:\\n if os.path.exists(p) and os.path.getsize(p)>0:\\n  with open(p) as f: d=json.load(f) or {}\\nexcept Exception: d={}\"); d.setdefault('wifi_mode','ap'); d.setdefault('ap_ssid','PiCamera'); d.setdefault('ap_password','picamera123'); t=p+'.tmp'; json.dump(d,open(t,'w'),indent=2); os.replace(t,p)"
ExecStart=/usr/bin/python3 /home/pi/api_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
