# home 3/camera-service.service
[Unit]
Description=Light Detection Camera Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
RuntimeDirectory=picamera
RuntimeDirectoryMode=0755
ExecStartPre=+/bin/sh -c 'echo performance > /sys/devices/system/cpu/cpu0/cpufreq/scaling_governor 2>/dev/null || true'
ExecStartPre=+/bin/sh -c 'echo deadline > /sys/block/mmcblk0/queue/scheduler 2>/dev/null || true'
ExecStartPre=+/bin/sh -c 'echo 0 > /proc/sys/vm/swappiness 2>/dev/null || true'
ExecStart=/usr/bin/python3 /home/pi/camera_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
