# home 3/picamera-time-sync.service
[Unit]
Description=Pi Camera: restore hardware clock and enable NTP
After=local-fs.target
Before=camera-service.service api-server.service
DefaultDependencies=yes

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
ExecStart=/home/pi/setup_zero2w_system.sh --time-sync-only

[Install]
WantedBy=multi-user.target
